# config_backup/systemd/user/change_refresh_rate.timer
[Unit]
Description=Run change_refresh_rate.service every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=default.target

